# com.marketdata.daily-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marketdata.daily-update</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/junlim/projects/market_data/daily_update.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>TimeZone</key>
        <string>America/New_York</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/junlim/projects/market_data/logs/daily_update.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/junlim/projects/market_data/logs/daily_update_error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/junlim/projects/market_data</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist> 